# launchd/com.nexus.cycle-tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nexus.cycle-tick</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/.cargo/bin/nexus</string>
        <string>singularity-tick</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ghost/Dev/nexus</string>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/Users/ghost/Library/Logs/nexus/cycle-tick.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ghost/Library/Logs/nexus/cycle-tick.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
